# common/systemd/user/battery-check.service
[Unit]
Description=Battery check notification service

[Service]
ExecStart=%h/bin/battery-check
